# deployment/168vape.service
[Unit]
Description=168VAPE Next.js Application
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/168vape
Environment=NODE_ENV=production
Environment=PORT=3000
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=168vape

[Install]
WantedBy=multi-user.target

